# rootfs/etc/services.d/claude-agent/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Claude Code agent stopping, unregistering from c3po..."

CREDS_FILE="/root/.claude/c3po-credentials.json"
if [ -f "$CREDS_FILE" ]; then
    coord_url=$(jq -r '.coordinator_url' "$CREDS_FILE" 2>/dev/null || echo "")
    cred_token=$(jq -r '.api_token' "$CREDS_FILE" 2>/dev/null || echo "")
    if [ -n "$coord_url" ] && [ "$coord_url" != "null" ] && \
       [ -n "$cred_token" ] && [ "$cred_token" != "null" ]; then
        curl -s -o /dev/null --max-time 5 -X POST \
            -H "Authorization: Bearer $cred_token" \
            -H "Content-Type: application/json" \
            -d '{}' \
            "$coord_url/agent/api/unregister" || true
        bashio::log.info "Unregistered from c3po"
    fi
fi

bashio::log.info "Claude Code agent stopped"
