# rootfs/etc/services.d/claude-agent/run
#!/usr/bin/with-contenv bashio

# Read config
TOKEN=$(bashio::config 'oauth_token')
WORK_DIR=$(bashio::config 'work_dir')
PROJECT_NAME=$(bashio::config 'project_name')
MACHINE_NAME=$(bashio::config 'machine_name')

# Export env for Claude and c3po
export CLAUDE_CODE_OAUTH_TOKEN="$TOKEN"
export CLAUDE_PROJECT_NAME="$PROJECT_NAME"
export C3PO_MACHINE_NAME="$MACHINE_NAME"
export C3PO_PROJECT_NAME="$PROJECT_NAME"
export C3PO_KEEP_REGISTERED="1"
export TMPDIR=/root/.claude/tmp
export PYTHONUNBUFFERED=1
export HOME=/root
# Allow --dangerously-skip-permissions to run as root (container is a sandbox)
export IS_SANDBOX=1

# Export custom env vars
if bashio::config.has_value 'env_vars'; then
    for entry in $(bashio::config 'env_vars'); do
        export "$entry"
        bashio::log.info "env: $entry"
    done
fi

# Build model flag
MODEL_FLAG=""
if bashio::config.has_value 'model'; then
    MODEL_FLAG="--model $(bashio::config 'model')"
fi

# Build prompt flag
PROMPT_FLAG=""
if bashio::config.has_value 'prompt'; then
    PROMPT_FLAG="$(bashio::config 'prompt')"
fi

bashio::log.info "Starting claude-watcher in ${WORK_DIR}..."
cd "$WORK_DIR"

exec python3 /usr/local/bin/claude-watcher.py \
    --work-dir "$WORK_DIR" \
    --model-flag "$MODEL_FLAG" \
    --prompt "$PROMPT_FLAG" \
    --session-dir /data/sessions \
    --max-sessions 50
