# rootfs/etc/services.d/claude-agent/run
#!/usr/bin/with-contenv bashio

ADDON_VERSION="0.1.5"
bashio::log.info "Starting Claude Code agent v${ADDON_VERSION}..."
bashio::log.info "Claude Code version: $(claude --version 2>&1 || echo 'unknown')"

# Validate OAuth token
if ! bashio::config.has_value 'oauth_token'; then
    bashio::log.fatal "OAuth token required in add-on configuration"
    exit 1
fi

TOKEN=$(bashio::config 'oauth_token')
if [[ ! "$TOKEN" =~ ^sk-ant- && ! "$TOKEN" =~ ^sk-at- ]]; then
    bashio::log.fatal "Invalid token format (must start with sk-ant- or sk-at-)"
    exit 1
fi

# Validate c3po configuration
if ! bashio::config.has_value 'c3po_coordinator_url'; then
    bashio::log.fatal "c3po coordinator URL required"
    exit 1
fi

C3PO_URL=$(bashio::config 'c3po_coordinator_url')
WORK_DIR=$(bashio::config 'work_dir')
PROJECT_NAME=$(bashio::config 'project_name')
MACHINE_NAME=$(bashio::config 'machine_name')

# Initialize .claude.json if it doesn't exist (for MCP config persistence)
if [ ! -f /home/node/.claude.json ]; then
    echo '{}' > /home/node/.claude.json
    chown node:node /home/node/.claude.json
fi

# Set up c3po plugin on first run
# After enrollment, c3po stores credentials in /home/node/.claude/c3po-credentials.json
# which persists via symlink to /data/claude/c3po-credentials.json
if [ ! -f /data/.c3po-setup-complete ]; then
    bashio::log.info "Setting up c3po plugin (first run)..."

    # Admin token only required for initial enrollment
    if ! bashio::config.has_value 'c3po_admin_token'; then
        bashio::log.fatal "c3po admin token required for first-time setup"
        exit 1
    fi

    C3PO_TOKEN=$(bashio::config 'c3po_admin_token')

    # TMPDIR must be on the same filesystem as ~/.claude for atomic renames
    mkdir -p /home/node/.claude/tmp
    chown node:node /home/node/.claude/tmp
    export TMPDIR=/home/node/.claude/tmp
    export CLAUDE_CODE_OAUTH_TOKEN="$TOKEN"

    # Debug: filesystem layout
    bashio::log.info "Debug: symlink target: $(readlink -f /home/node/.claude)"
    bashio::log.info "Debug: df output:"
    df -h /home/node /home/node/.claude /tmp 2>&1 | while read line; do bashio::log.info "  $line"; done
    bashio::log.info "Debug: ls -la /home/node/.claude:"
    ls -la /home/node/.claude 2>&1 | while read line; do bashio::log.info "  $line"; done

    # Add marketplace (as node user)
    bashio::log.info "Adding c3po marketplace..."
    if s6-setuidgid node env TMPDIR="$TMPDIR" CLAUDE_CODE_OAUTH_TOKEN="$TOKEN" \
        claude plugin marketplace add michaelansel/claude-code-plugins 2>&1 | while read line; do bashio::log.info "  marketplace: $line"; done; then
        bashio::log.info "Marketplace add succeeded"
    else
        bashio::log.warning "Marketplace add failed (exit $?), trying update..."
        s6-setuidgid node env TMPDIR="$TMPDIR" CLAUDE_CODE_OAUTH_TOKEN="$TOKEN" \
            claude plugin marketplace update michaelansel 2>&1 | while read line; do bashio::log.info "  marketplace: $line"; done
    fi

    # Install plugin (as node user)
    bashio::log.info "Installing c3po plugin..."
    if s6-setuidgid node env TMPDIR="$TMPDIR" CLAUDE_CODE_OAUTH_TOKEN="$TOKEN" \
        claude plugin install c3po@michaelansel 2>&1 | while read line; do bashio::log.info "  plugin: $line"; done; then
        bashio::log.info "Plugin install succeeded"
    else
        bashio::log.warning "Plugin install failed, trying update..."
        s6-setuidgid node env TMPDIR="$TMPDIR" CLAUDE_CODE_OAUTH_TOKEN="$TOKEN" \
            claude plugin update c3po@michaelansel 2>&1 | while read line; do bashio::log.info "  plugin: $line"; done
    fi

    # Debug: check what was created
    bashio::log.info "Debug: plugins directory after install:"
    ls -laR /home/node/.claude/plugins/ 2>&1 | while read line; do bashio::log.info "  $line"; done
    bashio::log.info "Debug: end plugins listing"

    # Find and run enrollment script
    SETUP_PY=$(find /home/node/.claude/plugins -path '*/c3po*/setup.py' -print -quit 2>&1)
    if [ -z "$SETUP_PY" ]; then
        bashio::log.fatal "Could not find c3po setup.py in /home/node/.claude/plugins"
        bashio::log.fatal "Plugin installation likely failed due to cross-filesystem issue"
        exit 1
    fi

    bashio::log.info "Running enrollment from $SETUP_PY..."
    s6-setuidgid node env CLAUDE_CODE_OAUTH_TOKEN="$TOKEN" \
        python3 "$SETUP_PY" --enroll "$C3PO_URL" "$C3PO_TOKEN" \
            --machine "$MACHINE_NAME" \
            --pattern "ha/*" 2>&1 | while read line; do bashio::log.info "  enroll: $line"; done

    # Verify credentials were created
    if [ ! -f /home/node/.claude/c3po-credentials.json ]; then
        bashio::log.fatal "c3po enrollment failed - no credentials file created"
        exit 1
    fi

    touch /data/.c3po-setup-complete
    bashio::log.info "c3po plugin setup complete"
    bashio::log.info "Machine credentials stored in /data/claude/c3po-credentials.json"
    bashio::log.info "You can now remove c3po_admin_token from add-on config (optional)"
else
    bashio::log.info "c3po already enrolled, using stored credentials"

    # Verify credentials still exist
    if [ ! -f /home/node/.claude/c3po-credentials.json ]; then
        bashio::log.fatal "c3po credentials missing! Remove /data/.c3po-setup-complete and restart to re-enroll"
        exit 1
    fi
fi

# Export token and project name for claude
export CLAUDE_CODE_OAUTH_TOKEN="$TOKEN"
export CLAUDE_PROJECT_NAME="$PROJECT_NAME"
export TMPDIR=/home/node/.claude/tmp

# Run claude with /c3po auto as node user
bashio::log.info "Starting /c3po auto in ${WORK_DIR}..."
cd "$WORK_DIR"
exec s6-setuidgid node claude --dangerously-skip-permissions -p "/c3po auto"
